# scripts/com.user.econ_digest_haiku.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.user.econ_digest_haiku.plist — launchd設定ファイル（Claude Haiku版）
  設置先: ~/Library/LaunchAgents/com.user.econ_digest_haiku.plist
  登録:   launchctl load ~/Library/LaunchAgents/com.user.econ_digest_haiku.plist
  手動実行: launchctl start com.user.econ_digest_haiku
  解除:   launchctl unload ~/Library/LaunchAgents/com.user.econ_digest_haiku.plist

  実行タイミング: 毎週金曜 13:00（Ollama版の09:00より4時間後）
  Ollama版（09:00）が先に実行されてから、Haikuで同じ週の記事を生成し、比較ページを作成する。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.econ_digest_haiku</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/masahiro/projects/econ_digest/scripts/run_econ_haiku.sh</string>
  </array>

  <!-- 実行スケジュール: 毎週金曜 13:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>5</integer>
    <key>Hour</key>
    <integer>13</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/masahiro/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/anaconda3/bin</string>
    <key>HOME</key>
    <string>/Users/masahiro</string>
    <key>LANG</key>
    <string>ja_JP.UTF-8</string>
    <!-- ANTHROPIC_API_KEY は ~/.anthropic_env から読み込まれる（シェルスクリプト内で処理） -->
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/masahiro/projects/econ_digest/econ_digest_haiku.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/masahiro/projects/econ_digest/econ_digest_haiku.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>WorkingDirectory</key>
  <string>/Users/masahiro/projects/econ_digest</string>
</dict>
</plist>
